# packages/filewave-monitor-v13/filewave-monitor-v13-1.0.5.tar.gz/filewave-monitor-v13-1.0.5/monitor/config/monitor-v13.conf
[program:apache_exporter]
command=/usr/local/sbin/apache_exporter --telemetry.address=:21108 --scrape_uri=https://localhost:20443/server-status?auto --insecure
autostart=true
priority=80
stdout_logfile=/usr/local/filewave/log/apache_exporter.out.log
stderr_logfile=/usr/local/filewave/log/apache_exporter.err.log

[program:mtail]
command=/usr/local/sbin/mtail --port 21110 --progs /usr/local/etc/filewave/mtail/progs --logs "/usr/local/filewave/log/*.log" --logs "/usr/local/filewave/fwxserver/DB/pg_data/pg_log/*.log" --logs "/usr/local/filewave/apache/logs/*_log*" --logs "/private/var/log/fw*.log"
autostart=true
priority=80
stdout_logfile=/usr/local/filewave/log/mtail.out.log
stderr_logfile=/usr/local/filewave/log/mtail.err.log

[program:prometheus]
command=/usr/local/sbin/prometheus --config.file=/usr/local/etc/filewave/prometheus/prometheus.yml --web.listen-address=:21090 --storage.tsdb.path=/usr/local/filewave/instrumentation_data/prometheus --web.console.templates=/usr/local/filewave/prometheus/consoles --web.console.libraries=/usr/local/filewave/prometheus/console_libraries
autostart=true
priority=90
stdout_logfile=/usr/local/filewave/log/prometheus.out.log
stderr_logfile=/usr/local/filewave/log/prometheus.err.log

[program:postgres_exporter]
command=/usr/local/sbin/postgres_exporter --web.listen-address=127.0.0.1:21109
autostart=true
priority=80
stdout_logfile=/usr/local/filewave/log/postgres_exporter.out.log
stderr_logfile=/usr/local/filewave/log/postgres_exporter.err.log
environment=DATA_SOURCE_NAME="user=postgres host=/tmp port=9432 sslmode=disable"

[program:node_exporter]
command=/usr/local/sbin/node_exporter --collector.supervisord --web.listen-address=127.0.0.1:21100
autostart=true
priority=80
stdout_logfile=/usr/local/filewave/log/node_exporter.out.log
stderr_logfile=/usr/local/filewave/log/node_exporter.err.log
